# deploy/vast/supervisor.conf.example
[program:srs]
command=/workspace/LiveTalking/deploy/vast/srs.sh
directory=/workspace/LiveTalking
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/portal/srs.log
stderr_logfile=/var/log/portal/srs-error.log

[program:runtime-manager]
command=/workspace/LiveTalking/deploy/vast/runtime-manager.sh
directory=/workspace/LiveTalking
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/portal/runtime-manager.log
stderr_logfile=/var/log/portal/runtime-manager-error.log
